# deploy/com.planningforms.redditbot.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd job for 24/7 operation. Keeps the bot running as a daemon
  (`--loop`), which starts human-like sessions at RANDOMIZED intervals — nothing
  on a fixed schedule. KeepAlive restarts it if it ever exits.

  SETUP:
    1. Build once:  npm run build   (the daemon below runs the compiled code)
    2. cp deploy/com.planningforms.redditbot.daemon.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.planningforms.redditbot.daemon.plist
    4. To stop: launchctl unload ~/Library/LaunchAgents/com.planningforms.redditbot.daemon.plist

  Use this OR the cron/single-run job — not both.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.planningforms.redditbot.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/etech/.local/node-22/bin/node</string>
    <string>/Users/etech/Desktop/zach-bot/dist/index.js</string>
    <string>--loop</string>
    <string>--live</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/etech/Desktop/zach-bot</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/etech/.local/node-22/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/etech/Desktop/zach-bot/logs/daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/etech/Desktop/zach-bot/logs/daemon.err.log</string>
</dict>
</plist>
